# resources/harbor/assets/flanneld/flanneld.path
[Path]
PathExists=/run/flannel/subnet.env
Unit=docker.service

[Install]
WantedBy=multi-user.target
